# deploy/systemd/miemie-worker.service
[Unit]
Description=Miemie Podcast Worker
After=network.target miemie-api.service

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
EnvironmentFile=__ROOT_DIR__/.env
ExecStart=/usr/bin/env APP_ENV=production PYTHONUNBUFFERED=1 __PYTHON_BIN__ -m miemie_podcast.worker.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
