# deploy/systemd/miemie-web.service
[Unit]
Description=Miemie Podcast Web
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/miemie-podcast/apps/web
EnvironmentFile=/opt/miemie-podcast/.env
Environment=APP_ENV=production
Environment=NODE_OPTIONS=--no-experimental-webstorage
ExecStart=/usr/bin/npx next start --hostname ${APP_HOST} --port ${APP_PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
