# deploy/systemd/miemie-api.service
[Unit]
Description=Miemie Podcast API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/miemie-podcast
EnvironmentFile=/opt/miemie-podcast/.env
Environment=APP_ENV=production
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/miemie-podcast/.venv/bin/python -m miemie_podcast.api.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
