# deploy/systemd/miemie-web.service
[Unit]
Description=Miemie Podcast Web
After=network.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__/apps/web
EnvironmentFile=__ROOT_DIR__/.env
Environment=APP_ENV=production
Environment=NODE_OPTIONS=--no-experimental-webstorage
ExecStart=__NPX_BIN__ next start --hostname ${APP_HOST} --port ${APP_PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
