# deploy/systemd/miemie-api.service
[Unit]
Description=Miemie Podcast API
After=network.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
EnvironmentFile=__ROOT_DIR__/.env
Environment=APP_ENV=production
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON_BIN__ -m miemie_podcast.api.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
